# server/support-files/release/supervisor/celery.conf
[program:celery]
; TODO: keep a single default Celery worker in this rollout. If long-running
; controller install tasks begin to contend with short tasks, split this into
; dedicated short/long workers and queues in a later optimization pass.
command=celery -A apps.core.celery worker --loglevel=info --pool threads --concurrency=%(ENV_CELERY_CONCURRENCY)s
directory=/apps
autostart=true
priority=10
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
